# packaging/scripts/postinstall.sh
systemctl daemon-reload
systemctl enable karbonara.service
systemctl restart karbonara.service
