# packaging/scripts/preremove.sh
systemctl stop karbonara.service
systemctl disable karbonara.service
